# ansible/roles/neon-agent/templates/multica-daemon.service.j2
[Unit]
Description=Multica daemon (cloud mode)
After=network-online.target
Wants=network-online.target

[Service]
User=neonuser
ExecStart=/usr/local/bin/multica daemon start
Restart=always
RestartSec=10
EnvironmentFile=/etc/neon.env

[Install]
WantedBy=multi-user.target
